# remote-root/etc/systemd/system/mihomo-manager.service
[Unit]
Description=Mihomo Web Manager
After=network.target
[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /etc/mihomo/manager/app.py
Restart=always
RestartSec=3
[Install]
WantedBy=multi-user.target
